# .docker/apache/production/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
file=/var/run/supervisor/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:httpd]
command=apache2-foreground
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/cron -l 2 -f
autostart=true
autorestart=true

[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --timeout=600 --memory=1024 --verbose
autostart=true
autorestart=true
user=root
numprocs=1 ; usually 3 or more would be better
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-worker.logs
